# recipes-support/boyd/files/boyd.service
[Unit]
Description=Boyd
After=local-fs.target daylited.service

[Service]
ExecStart=/usr/bin/boyd

[Install]
WantedBy=multi-user.target
